# packaging/Debian/debian/samba.conffiles
/etc/cron.daily/samba
/etc/init.d/samba
/etc/logrotate.d/samba
