# packaging/Debian/debian/winbind.conffiles
/etc/init.d/winbind
/etc/logrotate.d/winbind
