# packaging/Debian/unstable/samba.postrm
#!/bin/sh

if [ "$1" = purge ]; then
	update-rc.d samba remove >/dev/null

	# Remove WINS.DAT, BROWSE.DAT and lock information file
	rm -Rf /var/samba/

	# Remove any files in the old and obsolete /var/lock/samba directory
	rm -Rf /var/lock/samba/

	# Remove files left in /etc/samba/
	rm -Rf /etc/samba/debian_config
	rm -Rf /etc/samba/MACHINE.SID

	# Remove log files
	rm -f /var/log/[ns]mb*

	# Remove NetBIOS entries from /etc/inetd.conf
	update-inetd --remove netbios-ssn
	update-inetd --remove netbios-ns
else
	# Not purging, do not remove NetBIOS entries from /etc/inetd.conf
	update-inetd --disable netbios-ssn
	update-inetd --disable netbios-ns
fi
