# packaging/systemd/nmb.service.in
[Unit]
Description=Samba NMB Daemon
Wants=network-online.target
After=syslog.target network.target network-online.target

[Service]
Type=notify
NotifyAccess=all
PIDFile=@PIDDIR@/nmbd.pid
EnvironmentFile=-@SYSCONFDIR@/sysconfig/samba
ExecStart=@SBINDIR@/nmbd --foreground --no-process-group $NMBDOPTIONS
ExecReload=/usr/bin/kill -HUP $MAINPID
LimitCORE=infinity
@systemd_nmb_extra@

[Install]
WantedBy=multi-user.target
