# packaging/systemd/samba.service.in
[Unit]
Description=Samba AD Daemon
Wants=network-online.target
After=syslog.target network.target network-online.target

[Service]
Type=notify
NotifyAccess=all
PIDFile=@PIDDIR@/samba.pid
LimitNOFILE=16384
EnvironmentFile=-@SYSCONFDIR@/sysconfig/samba
ExecStart=@SBINDIR@/samba --foreground --no-process-group $SAMBAOPTIONS
ExecReload=/usr/bin/kill -HUP $MAINPID
@systemd_samba_extra@

[Install]
WantedBy=multi-user.target
